# systemd/rdp-honeypot-logs.service
[Unit]
Description=Process RDP honeypot logs (publish public/private)
After=rdp-honeypot.service

[Service]
Type=oneshot
EnvironmentFile=-@@INSTALL_DIR@@/.env
Environment=HONEYPOT_LOG_DIR=@@INSTALL_DIR@@/data/logs
ExecStart=/usr/bin/python3 @@INSTALL_DIR@@/scripts/log_processor.py
